# Homeworks/Module 2. Linux Basics/2.7 Bash best practices/10/watch.service
[Unit]
Description=Cahching script
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/vaclav
ExecStart=/home/vaclav/script.sh
Environment="PATH=/usr/bin:/bin"